# web/golf-pool-web.example.service
# Example systemd unit — copy to /etc/systemd/system/golf-pool-web.service and adjust paths.
#
#   sudo cp web/golf-pool-web.example.service /etc/systemd/system/golf-pool-web.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now golf-pool-web
#
# Nginx Proxy Manager: proxy to http://127.0.0.1:5000 (or the port you set below).

[Unit]
Description=Golf pool Flask site (gunicorn)
After=network.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_GROUP
WorkingDirectory=/absolute/path/to/golf
Environment=GOLF_POOL_ROOT=/absolute/path/to/golf
Environment=PATH=/absolute/path/to/golf/.venv/bin
ExecStart=/absolute/path/to/golf/.venv/bin/gunicorn -w 2 -b 127.0.0.1:5000 wsgi:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
